# launchd/com.wallpaper-aerials-fix.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wallpaper-aerials-fix</string>

  <key>ProgramArguments</key>
  <array>
    <!-- Update this path to your built binary -->
    <string>/ABSOLUTE/PATH/TO/wallpaper-aerials-fix</string>
    <string>--event</string>
    <string>unlock</string>
    <!-- Uncomment the next line if you have external monitors -->
    <!-- <string>--wait-for-displays</string> -->
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/dev/null</string>

  <key>StandardErrorPath</key>
  <string>/dev/null</string>
</dict>
</plist>
